# templates/seaweedfs-s3.service.j2
[Unit]
Description=SeaweedFS S3 Gateway
After=network-online.target seaweedfs-filer.service
Wants=network-online.target

[Service]
Type=simple
User={{ seaweedfs_user }}
Group={{ seaweedfs_group }}
ExecStart={{ seaweedfs_bin }} s3 \
  -ip.bind={{ seaweedfs_s3_bind_ip }} \
  -port={{ seaweedfs_s3_port }} \
  -port.grpc={{ seaweedfs_s3_grpc_port }} \
  -filer={% for h in groups.get('filers', []) %}{{ hostvars[h]['ansible_host'] }}:{{ seaweedfs_filer_port }}{% if not loop.last %},{% endif %}{% endfor %} \
  -config=/etc/seaweedfs/s3.conf
WorkingDirectory=/etc/seaweedfs
Restart=always
RestartSec=5
LimitNOFILE=1048576
StandardOutput=append:{{ seaweedfs_log_dir }}/s3.log
StandardError=append:{{ seaweedfs_log_dir }}/s3.err.log

[Install]
WantedBy=multi-user.target